# system/homeserver-backup.timer
[Unit]
Description=HOMESERVER Backup Timer - Automated backup scheduling
Requires=homeserver-backup.service

[Timer]
# Daily backup at 2 AM
OnCalendar=*-*-* 02:00:00
# Weekly backup on Sunday at 3 AM
OnCalendar=weekly
OnCalendar=*-*-* 03:00:00
# Monthly backup on 1st at 4 AM
OnCalendar=monthly
OnCalendar=*-*-01 04:00:00
# Yearly backup on January 1st at 5 AM
OnCalendar=yearly
OnCalendar=*-01-01 05:00:00

# Randomize start time to avoid thundering herd
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
