# backend/src/config/homeserver-backup.service
[Unit]
Description=HOMESERVER Backup Service
Documentation=https://homeserver.com/docs/backup
After=network.target homeserver.service
Wants=homeserver.service

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/homeserver/premium/backupTab
ExecStart=/usr/bin/python3 /var/www/homeserver/premium/backupTab/backend/src/service/backup_service.py --backup
StandardOutput=journal
StandardError=journal
SyslogIdentifier=homeserver-backup

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/homeserver/premium/backupTab/backend/src /var/log/homeserver
ReadOnlyPaths=/var/www/homeserver/premium/backupTab/src /var/lib/gogs /etc/homeserver


[Install]
WantedBy=multi-user.target
